# roles/install_monitoring/files/ffmyk-influx/ffmyk-influx.service
[Unit]
 
[Service]
ExecStart=/opt/ffmyk-influx/daemon.sh
Restart=always
RestartSec=5
 
[Install]
WantedBy=multi-user.target

