# roles/install_wglinkstat/templates/wg-linkstat@.service
[Unit]
Description=collect statistics on wireguard interface %I
After=network.target

[Service]
ExecStart=/usr/local/bin/wg_linkstat.sh "%i"
Environment=INFLUX_USER={{ hostvars['ffww-map'].influx.username }}
Environment=INFLUX_PASSWORD={{ hostvars['ffww-map'].influx.password }}
Environment=INFLUX_BASE_URL={{ hostvars['ffww-map'].influx.baseUrl }}
Environment=INFLUX_DB={{ hostvars['ffww-map'].influx.db }}
WorkingDirectory=/run/
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target
